# deploy/binance-bot-backup.timer
[Unit]
Description=Daily backup of bot.db

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
